# pb.ini
[supervisord]
nodaemon=true

[program:php-fpm]
priority=1
command=/usr/sbin/php-fpm7 -R
autorestart=true

[program:nginx]
priority=2
command=/usr/sbin/nginx
autorestart=true